# error-to-slack.service
[Unit]
Description=Rust Error to Slack

[Service]
ExecStart=/usr/sbin/rust-errors-to-slack

[Install]
WantedBy=multi-user.target
